# roles/launchd/templates/com.ahmedelgabri.isync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple/DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.ahmedelgabri.isync</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>exec {{ '~/.config/zsh/bin/mail-sync' | expanduser }}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>StandardOutPath</key>
    <string>/tmp/ahmed.isync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ahmed.isync.err.log</string>
  </dict>
</plist>
